# deploy/dartplayr.service
[Unit]
Description=DartPlayr Django app
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/dartplayr
EnvironmentFile=/etc/dartplayr.env
ExecStart=/var/www/dartplayr/venv/bin/gunicorn dartplayr.wsgi:application --workers 3 --bind unix:/run/dartplayr/gunicorn.sock
Restart=always
RestartSec=5
RuntimeDirectory=dartplayr

[Install]
WantedBy=multi-user.target
